# supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=sandbox
pidfile=~/.openclaw/logs/supervisord.pid
logfile=~/.openclaw/logs/supervisord.log
[inet_http_server]
port = 127.0.0.1:9001

[supervisorctl]
serverurl=http://127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:claw-worker]
command=python3 /home/sandbox/.openclaw/extensions/galaxyos/scripts/claw_worker.py --daemon
directory=/home/sandbox
autostart=true
autorestart=true
startsecs=2
startretries=10
stopwaitsecs=5
killasgroup=true
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1
environment=OPENCLAW_WORKSPACE=/home/sandbox/.openclaw/workspace,PYTHONPATH=/home/sandbox/.openclaw/workspace/repo/lib/python3.12/site-packages:/home/workspace/.local/lib/python3.12/site-packages/:/home/sandbox/.local/lib/python3.12/site-packages/,MKL_NUM_THREADS=1,OMP_NUM_THREADS=1

[program:openclaw-gateway]
command=openclaw gateway
directory=/home/sandbox
autostart=true
autorestart=true
startsecs=3
startretries=10
stopwaitsecs=2          ; 等待10秒优雅关闭
killasgroup=true         ; 杀死整个进程组
stopsignal=TERM          ; 使用 SIGTERM 信号
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_COMPILE_CACHE=/tmp/openclaw-compile-cache,OPENCLAW_NO_RESPAWN=1,PYTHONPATH=/home/sandbox/.openclaw/workspace/repo/lib/python3.12/site-packages:/home/workspace/.local/lib/python3.12/site-packages/:/home/sandbox/.local/lib/python3.12/site-packages/,MKL_NUM_THREADS=1,OMP_NUM_THREADS=1,LD_PRELOAD=/home/sandbox/.local/lib/libmkl_wrapper.so,LD_LIBRARY_PATH=/home/sandbox/.local/mkl/lib


[program:apscheduler]
command=python3 /home/sandbox/scheduler.py
directory=/home/sandbox
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=10          ; 等待10秒优雅关闭
killasgroup=true         ; 杀死整个进程组
stopsignal=TERM          ; 使用 SIGTERM 信号
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_COMPILE_CACHE=/tmp/openclaw-compile-cache,OPENCLAW_NO_RESPAWN=1


[program:self_heal]
command=python /opt/bin/self_heal.py
directory=/home/sandbox
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=10          ; 等待10秒优雅关闭
killasgroup=true         ; 杀死整个进程组
stopsignal=TERM          ; 使用 SIGTERM 信号
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_COMPILE_CACHE=/tmp/openclaw-compile-cache,OPENCLAW_NO_RESPAWN=1


[program:watch_paired]
command=python /opt/bin/watch_paired.py
directory=/home/sandbox
autostart=true
autorestart=false
startsecs=2
killasgroup=true         ; 杀死整个进程组
stopsignal=TERM          ; 使用 SIGTERM 信号
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_COMPILE_CACHE=/tmp/openclaw-compile-cache,OPENCLAW_NO_RESPAWN=1